# config/moxnas.service
[Unit]
Description=MoxNAS Web Interface
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/moxnas
Environment=PATH=/opt/moxnas/venv/bin
ExecStart=/opt/moxnas/venv/bin/python /opt/moxnas/start_moxnas.py --production
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target